# crontab
0,14,29,44 * * * * /home/pablo/Documents/git-projects/henry/labs/Streamlit_datathon/datathon.sh
